# systemd/pytime/pytime@.service
[Unit]
Description=Pytime Service    %i

[Service]
Type=oneshot
ExecStart=/%I
